# sectyra/codespace_data/crontab
0 * * * * root /bin/bash /home/codespace/codespace_data/keep_alive.sh 2>/home/codespace/error.log
